# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
ipset -N youtubev4 hash:net
ipset -N youtubev6 hash:net family inet6
service dnsmasq restart
mwan3 restart
exit 0
